# mssql/cron.d/mssql
* * * * * root /root/git/docker/mssql/backup.sh mssql@data2.regatta-hd.de
